# llm_news.service
[Unit]
Description=LLM Haber Ozeti — Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bllsrky
WorkingDirectory=/home/bllsrky/srky
ExecStart=/usr/bin/python3 /home/bllsrky/srky/run.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
